# support/debian/prerm
systemctl stop getnewip
systemctl disable getnewip
